# deploy/supervisor/thriive.conf
; Thriive IMS — Supervisor process group.
; Placeholders __APP_DIR__ / __APP_USER__ / __GUNICORN_WORKERS__ are substituted
; by deploy/setup.sh (workers default to 2*nproc+1, capped at 24).
; Long-running processes: Gunicorn, two Celery workers (compute + io), Celery beat.
;
;   supervisorctl status thriive:
;   supervisorctl restart thriive:

[group:thriive]
programs=web,celery_compute,celery_io,beat

[program:web]
; gthread workers: N processes x 4 threads of concurrency, so a few slow
; report/export requests can't occupy every worker. max-requests recycles
; workers to cap slow memory growth.
command=__APP_DIR__/backend/venv/bin/gunicorn config.wsgi:application
    --workers __GUNICORN_WORKERS__ --worker-class gthread --threads 4
    --max-requests 1000 --max-requests-jitter 100
    --bind 127.0.0.1:8000 --timeout 120 --access-logfile -
directory=__APP_DIR__/backend
user=__APP_USER__
autostart=true
autorestart=true
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/log/thriive/web.log
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"

[program:celery_compute]
; Payouts + achievements get a dedicated worker so month-close compute is never
; queued behind report generation or bulk imports.
command=__APP_DIR__/backend/venv/bin/celery -A config worker
    --loglevel=info -Q payouts,achievements
    --concurrency=4 --time-limit=600 --soft-time-limit=540
directory=__APP_DIR__/backend
user=__APP_USER__
autostart=true
autorestart=true
stopwaitsecs=600
redirect_stderr=true
stdout_logfile=/var/log/thriive/celery-compute.log
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"

[program:celery_io]
; 'celery' is Celery's own default queue name — kept in the list as a safety net so a task
; that slips through unrouted still runs instead of queueing into silence.
command=__APP_DIR__/backend/venv/bin/celery -A config worker
    --loglevel=info -Q default,reports,celery
    --concurrency=4 --time-limit=600 --soft-time-limit=540
directory=__APP_DIR__/backend
user=__APP_USER__
autostart=true
autorestart=true
stopwaitsecs=600
redirect_stderr=true
stdout_logfile=/var/log/thriive/celery-io.log
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"

[program:beat]
command=__APP_DIR__/backend/venv/bin/celery -A config beat
    --loglevel=info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=__APP_DIR__/backend
user=__APP_USER__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/thriive/beat.log
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"
