# deploy/supervisor/thriive.conf
; Thriive IMS — Supervisor process group.
; Placeholders __APP_DIR__ / __APP_USER__ are substituted by deploy/setup.sh.
; Controls the three long-running processes: Gunicorn, the Celery worker, Celery beat.
;
;   supervisorctl status thriive:
;   supervisorctl restart thriive:

[group:thriive]
programs=web,celery,beat

[program:web]
command=__APP_DIR__/backend/venv/bin/gunicorn config.wsgi:application
    --workers 4 --bind 127.0.0.1:8000 --timeout 120 --access-logfile -
directory=__APP_DIR__/backend
user=__APP_USER__
autostart=true
autorestart=true
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/log/thriive/web.log
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"

[program:celery]
command=__APP_DIR__/backend/venv/bin/celery -A config worker
    --loglevel=info -Q default,achievements,payouts,reports
    --concurrency=4 --time-limit=600 --soft-time-limit=540
directory=__APP_DIR__/backend
user=__APP_USER__
autostart=true
autorestart=true
stopwaitsecs=600
redirect_stderr=true
stdout_logfile=/var/log/thriive/celery.log
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"

[program:beat]
command=__APP_DIR__/backend/venv/bin/celery -A config beat
    --loglevel=info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=__APP_DIR__/backend
user=__APP_USER__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/thriive/beat.log
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"
