# deploy/supervisor/thriive.conf
[group:thriive]
programs=web,celery,beat

[program:web]
command=/opt/thriive/backend/venv/bin/gunicorn config.wsgi:application -w 4 -b 127.0.0.1:8000 --timeout 120
directory=/opt/thriive/backend
user=thriive
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/thriive/web.log
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"

[program:celery]
command=/opt/thriive/backend/venv/bin/celery -A config worker --loglevel=info -Q default,achievements,payouts,reports --concurrency=4 --time-limit=600 --soft-time-limit=540
directory=/opt/thriive/backend
user=thriive
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/thriive/celery.log
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"

[program:beat]
command=/opt/thriive/backend/venv/bin/celery -A config beat --loglevel=info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/opt/thriive/backend
user=thriive
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/thriive/beat.log
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"
